# dashboard/services/home-relay/home-relay.service
[Unit]
Description=Home Relay Service (Kasa + WOL)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=/home/__USER__/dashboard/services/home-relay
ExecStart=/home/__USER__/dashboard/services/home-relay/.venv/bin/python relay.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
