# dashboard/services/home-relay/voice-control.service
[Unit]
Description=Voice Control Service
After=network.target home-relay.service
Wants=home-relay.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=/home/__USER__/dashboard/services/home-relay
ExecStart=/home/__USER__/.local/bin/uv run --group voice python voice_control.py
Restart=on-failure
RestartSec=5

# Audio device access
SupplementaryGroups=audio

[Install]
# Started via Settings UI, not at boot
# WantedBy=multi-user.target
